# monitor.timer
[Unit]
Description=Run monitor every 1 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=120

[Install]
WantedBy=timers.target

