# monitor.timer
[Unit]
Description=Run monitor every 1 min

[Timer]
OnBootSec=5min
OnUnitActiveSec=50

[Install]
WantedBy=timers.target

