# monitor.timer
[Unit]
Description=Run monitor every 1 min

[Timer]
OnBootSec=10
OnUnitActiveSec=60

[Install]
WantedBy=timers.target

