# monitor.timer
[Unit]
Description=Run monitor every 15 min

[Timer]
OnBootSec=15min
OnUnitActiveSec=900

[Install]
WantedBy=timers.target

