# monitor.service
[Unit]
Description="Service Monitor script"

[Service]
ExecStart=/usr/bin/timeout 2m /usr/bin/python3 /home/pi/espRelease/monitor.py
