# modemManager/resetModem.service
[Unit]
Description="Service Monitor script"

[Service]
ExecStart=/usr/bin/timeout 2m /etc/indra/resetModem.sh | ts 

